# overlays/mysql/etc/cron.daily/mysqloptimize
#!/bin/sh

mysqloptimize --all-databases > /dev/null
