# overlays/mysql/etc/cron.daily/mysqloptimize
#!/bin/sh

mysqloptimize --defaults-file=/etc/mysql/debian.cnf --all-databases
