# overlays/turnkey.d/update-motd/etc/cron.d/update-motd
# cron job for updating motd
*/10 * * * * root run-parts /etc/update-motd.d > /etc/motd


