# overlays/turnkey.d/webmin/etc/systemd/system/webmin.service.d/override.conf
[Unit]
After=network-online.target remote-fs.target stunnel4@webmin.service
Requires=stunnel4@webmin.service
BindsTo=stunnel4@webmin.service
