# overlays/pulseaudio-alsa/etc/udev/rules.d/86-pulseaudio.rules
KERNEL=="controlC[0-9]*", ACTION=="add", \
	RUN+="/sbin/start-stop-daemon --start --background --pidfile /var/run/pulseaudio/bogus --startas /etc/init.d/pulseaudio -- start %n"
